# infra/image/mkosi.skeleton/etc/systemd/system/tdx-launcher.service
# TDX VM Launcher Service
# Reads config from /mnt/config/config.json (verity) or /etc/easyenclave/config.json (legacy)
# and runs in control-plane or agent mode.

[Unit]
Description=TDX VM Launcher Service
After=data.mount mnt-config.mount docker.service network-online.target
Wants=docker.service network-online.target
# Don't hard-fail if config disk is missing (legacy boot)
After=mnt-config.mount
Wants=mnt-config.mount

[Service]
Type=simple
User=root
WorkingDirectory=/opt/launcher
ExecStart=/usr/bin/python3 /opt/launcher/launcher.py
Restart=on-failure
RestartSec=5
StandardOutput=journal+console
StandardError=journal+console

[Install]
WantedBy=multi-user.target
