# infra/image/mkosi.skeleton/etc/systemd/system/tdx-tsm-setup.service
# Setup ConfigFS-TSM for TDX attestation quote generation

[Unit]
Description=Setup ConfigFS-TSM for TDX attestation
After=local-fs.target

[Service]
Type=oneshot
ExecStart=/bin/bash -c 'modprobe configfs 2>/dev/null || true; mount -t configfs none /sys/kernel/config 2>/dev/null || true'
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
